# config/supervisor/asterisk.conf
[program:asterisk]
command=/usr/sbin/asterisk -g -f -U asterisk
startsecs=30
autorestart=true
user=root
group=root
